# spamhaushbl.conf
[program:spamhaushbl]
command=/etc/milter/spamhaushbl.py
process_name=spamhaushbl
priority=1
redirect_stderr=true
stdout_logfile=/var/log/spamhaushbl.log
